# publishing.ini
[program:publishing]
command=gunicorn publish_django.wsgi:application --bind localhost:8080 --log-file /var/log/publishing.log
directory=/var/publish_django/
stdout_logfile=/var/log/publishing.log
autostart=true
autorestart=true
startretries=3
